# etc/rc.local
#!/bin/sh -e

/etc/network-namespaces.d/netns_init.sh

exit 0
