# etc/systemd/system/iptables.service.d/override.conf
[Service]
ExecStart=/usr/local/bin/iptables-netns
